# com.dustpredictor.tashkent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dustpredictor.tashkent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/bulatdavlyatshin/Working/Projects/experiments/python/dustcast/tashkent_dust_predictor/run.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/bulatdavlyatshin/Working/Projects/experiments/python/dustcast/tashkent_dust_predictor</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/bulatdavlyatshin/Working/Projects/experiments/python/dustcast/tashkent_dust_predictor/logs/launchd.log</string>
    <key>StandardOutPath</key>
    <string>/Users/bulatdavlyatshin/Working/Projects/experiments/python/dustcast/tashkent_dust_predictor/logs/launchd.log</string>
</dict>
</plist>
